# deploy/palazzo.service
[Unit]
Description=palazzo — Qdrant memory palace MCP server (Streamable HTTP)
Documentation=https://github.com/calibrae/palazzo
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=palazzo
Group=palazzo
ExecStart=/opt/palazzo/bin/palazzo serve
Restart=on-failure
RestartSec=5

# palazzo reads Qdrant/Ollama URLs from env; override in /etc/palazzo/env
EnvironmentFile=-/etc/palazzo/env
Environment=RUST_LOG=palazzo=info

# Hardening — palazzo only needs outbound HTTP to Qdrant/Ollama and writes WAL.
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/palazzo
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
